# sys-apps/sshd-configurator/files/sshd-configurator.initd
#!/sbin/openrc-run
# Copyright 1999-2019 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/sshd-configurator"
command_args="${SSHD_INTERFACE}"

depend() {
	before sshd
	use logger
	provide sshd-configurator
}
